# meta/recipes-core/runit/runit-serialgetty/run
#!/bin/sh
exec /sbin/getty @BAUDRATE@ @TTY@ vt102
